# docker/build/supervisord-z-create-mysql-user.conf
[program:create_mysql_user]
command=/create_mysql_user.sh
startsecs = 0
autorestart = false