# docker/build/supervisord-z-create-mysql-user.conf
[program:create_mysql_user]
command=/create_mysql_user.sh
numprocs=1
autostart=true
autorestart=false