# docker/build/supervisord-server-symfony.conf
[program:symfony_server]
command=/create_mysql_user.sh
numprocs=1
autostart=true
autorestart=true